# deepform.service
[Unit]
Description=DeepForm AI Depth Interview
After=network.target

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/DeepForm
ExecStart=/usr/bin/node --experimental-sqlite --import file:///home/exedev/DeepForm/node_modules/tsx/dist/loader.mjs --require /home/exedev/DeepForm/node_modules/tsx/dist/preflight.cjs src/index.ts
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production
EnvironmentFile=/home/exedev/DeepForm/.env

[Install]
WantedBy=multi-user.target
