# deepform.service
[Unit]
Description=DeepForm AI Depth Interview
After=network.target

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/DeepForm
ExecStart=/home/exedev/DeepForm/node_modules/.bin/tsx src/index.ts
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
